# scripts/crontab
*/5 * * * * /scripts/sync-repos.sh >> /var/log/cron.log 2>&1
